# setup/sma-daemon.service
[Unit]
Description=SMA-NG Daemon - Media Conversion Webhook Server
Documentation=https://github.com/newdave/sma-ng
After=network.target
Wants=network.target
# When using a local PostgreSQL instance, uncomment the following two lines
# so the daemon starts after the database is ready:
#After=network.target postgresql.service
#Wants=postgresql.service

[Service]
Type=simple
User=nobody
Group=nogroup

WorkingDirectory=/opt/sma

# Load site-local environment overrides (optional — missing file is not an error)
EnvironmentFile=-/opt/sma/config/daemon.env

# Defaults — override any of these in config/daemon.env
Environment="SMA_DAEMON_HOST=0.0.0.0"
Environment="SMA_DAEMON_PORT=8585"
Environment="SMA_DAEMON_WORKERS=1"
Environment="SMA_CONFIG=/opt/sma/config/autoProcess.ini"
# Set this to use PostgreSQL instead of SQLite.
# When set, --db (SQLite) is automatically skipped.
# Format: postgresql://user:password@host:port/dbname
#Environment="SMA_DAEMON_DB_URL=postgresql://sma:password@localhost:5432/sma"

# Wrapper script selects --db-url (PostgreSQL) or --db (SQLite) based on
# whether SMA_DAEMON_DB_URL is set in the environment.
ExecStart=/opt/sma/setup/sma-daemon-start.sh

Restart=on-failure
RestartSec=10
TimeoutStopSec=infinity
KillMode=mixed

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sma/config /opt/sma/logs /transcodes /mnt

# Resource limits
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sma-daemon

[Install]
WantedBy=multi-user.target
